# bbs.conf
[program:bbs]
directory=/home/bbs/bbs
command=gunicorn wsgi -c gunicorn.config.py
autostart=true
autorestart=true